# roles/wifi_persistence/files/wifi-reinit.service
[Unit]
Description=Force Re-initialize WiFi on Boot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wifi-reinit.sh
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target