# roles/wifi_persistence/files/wifi-reinit.service
[Unit]
Description=Force Re-initialize WiFi on Boot
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi-reinit.sh
User=root

[Install]
WantedBy=multi-user.target